# cloud-init/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:cloud-init]
command=/tmp/cloud-init-run.sh
