# murano-agent/murano_agent_super.conf

[program:murano_agent]
command=/usr/local/bin/muranoagent
